# config/rules/99-gpu-power.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/home/swon/.scripts/gpu-autocheck.sh"
