# config/rules/99-gpu-power.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/usr/bin/sudo -u your_user_name /home/your_user_name/.scripts/gpu-autocheck.sh" 
